# ubuntu-php8.2/config/supervisord.conf
[supervisord]
nodaemon = true

[program:httpd]
command = apachectl -D FOREGROUND
